# srv1/mate-collector
#!/bin/sh /etc/rc.common
NAME=mate-collector
START=60
STOP=60

start() {
    echo "Starting pyMATE collector"
    /root/pymate/collector.py &
}

stop() {
    killall collector.py
}
